# products/coco/tika/supervisor/tika.conf
[program:tika]
command=/app/easysearch/data/coco/start-tika.sh
directory=/app/easysearch/data/coco
autostart=true
autorestart=true

# 日志配置
redirect_stderr=true
stdout_logfile=/app/easysearch/data/coco/supervisor.tika.log
stderr_logfile=NONE
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

# 启动参数调整（Java 启动可能需要更长时间）
startsecs=30
startretries=3

# 进程管理（现在直接管理 Java，不需要组杀）
stopasgroup=false 
killasgroup=false
stopsignal=TERM
stopwaitsecs=30

user=root